# configs/etc/crontab
# Generate sitemaps
0   5   *   *   *   www-data    /srv/LIBDb.Settings/sitemaps-generate.sh

# Every sunday at 5:00 backup LIBDB Files
0   5   *   *   0   www-data    /srv/backup-files-libdb.sh >/dev/null 2>&1

# Backup libdb databases
0   6   *   *   *   www-data    /srv/backup-db-libdb.sh >/dev/null 2>&1


